# LaunchAgents/com.ebv.feed.btc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ebv.feed.btc</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>-lc</string>
    <string>cd '$HOME/ebv_data_stack' && source .venv/bin/activate && python scripts/feed_coinbase_minute.py</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>$(HOME)/Library/Logs/ebv.feed.btc.out.log</string>
  <key>StandardErrorPath</key><string>$(HOME)/Library/Logs/ebv.feed.btc.err.log</string>
</dict>
</plist>
